# unix/packaging-misc/husarnet-license-refresh.timer
[Unit]
Description=Refresh Husarnet Self-hosted License Daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
